# package/deb/after-install.sh
#!/bin/bash
systemctl enable lookatch-agent.service
systemctl start lookatch-agent.service