# package/deb/after-install.sh
#!/bin/bash
update-rc.d lookatch-agent defaults
/etc/init.d/lookatch-agent start
