# package/deb/before-remove.sh
#!/bin/bash
/etc/init.d/lookatch-agent stop
update-rc.d -f lookatch-agent remove
