# package/rpm/before-remove.sh
#!/bin/bash
systemctl stop lookatch-agent.service
chkconfig --del lookatch-agent
