# package/rpm/before-remove.sh
#!/bin/bash
/etc/init.d/lookatch-agent stop
chkconfig --del lookatch-agent
